# rootfs/etc/services.d/nodered/run
#!/usr/bin/with-contenv sh

cd /usr/src/node-red || exit

echo "Starting node-red"
s6-setuidgid node-red npm start -- --userDir /data
